# utils/luci-theme-pberg/uci-defaults/luci-theme-pberg
#!/bin/sh
uci batch <<-EOF
	set luci.themes.pberg=/luci-static/pberg
        commit luci
EOF
	
