# utils/luci-theme-bootstrap-piraten/uci-defaults/luci-theme-bootstrap-piraten
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Bootstrap_Piraten=/luci-static/bootstrap-piraten
	set luci.main.mediaurlbase=/luci-static/bootstrap-piraten
	commit luci
EOF

