# utils/luci-app-ffwizard-berlin/root/etc/uci-defaults/wizarddefaults
#!/bin/sh

uci set freifunk.wizard=settings
uci commit
