# utils/luci-theme-piraten/uci-defaults/luci-theme-piraten
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Piraten=/luci-static/piraten
	set luci.main.mediaurlbase=/luci-static/piraten
        commit luci
EOF
	
