# uplinks/freifunk-berlin-uplink-no-tunnel/uci-defaults/freifunk-berlin-notunnel
#!/bin/sh

. /lib/functions/freifunk-berlin-network.sh

# always set correct masquerading, regardless of guard
uci set firewall.zone_ffuplink.masq=1
uci commit firewall

if [[ $(uci get ffberlin-uplink.preset.current) != "no-tunnel" ]]; then
  logger -t "ffuplink" "uplink-preset has been changed."
  # do not track preset when it was 'undefined', aka never configured
  if [[ $(uci get ffberlin-uplink.preset.current) != 'undefined' ]]; then
    uci rename ffberlin-uplink.preset.current=previous
  fi
  uci set ffberlin-uplink.preset.current="no-tunnel"
fi
# set set auth-type required for this uplink-type, e.g. for freifunk-wizard
uci set ffberlin-uplink.uplink=settings
uci set ffberlin-uplink.uplink.auth=none

uci commit ffberlin-uplink

. /lib/functions/guard.sh
guard "notunnel"

uci delete network.ffuplink_dev
uci set network.ffuplink_dev=device
uci set network.ffuplink_dev.type=veth
uci set network.ffuplink_dev.name=ffuplink
uci set network.ffuplink_dev.peer_name=ffuplink_wan
# add ffuplink_dev to br-wan if not there
uci set network.wan.ifname="$(uci get network.wan.ifname) ffuplink_wan"
uci commit network.ffuplink_dev

create_ffuplink
uci set network.ffuplink.proto=dhcp
uci commit network.ffuplink
