# defaults/freifunk-berlin-network-defaults/uci-defaults/freifunk-berlin-network-ip-route
#!/bin/sh

. /lib/functions/guard.sh
guard "iproute"

echo 120 ffuplink >>/etc/iproute2/rt_tables
